# deploy/systemd/claw-backup.service
# claw-backup.service — One-shot backup service triggered by claw-backup.timer

[Unit]
Description=Claw Orchestrator Daily Backup
After=network.target

[Service]
Type=oneshot
User=claw
WorkingDirectory=__REPO_DIR__
EnvironmentFile=-/etc/claw-orchestrator/env
ExecStart=__REPO_DIR__/deploy/scripts/backup.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claw-backup
